# calibre_web_automated/rootfs/etc/s6-overlay/s6-rc.d/ha-cwa-options/up
/etc/s6-overlay/s6-rc.d/ha-cwa-options/run
